# fs-overlay/fs/etc/systemd/system/picam-rtsp-server.service
[Unit]
Description=PiCam RTSP Server

[Service]
User=pi
ExecStart=/opt/stream/bin/picam-rtsp-server

[Install]
WantedBy=multi-user.target
